# zxd-bbs.conf
[program:zxd-bbs]
command=/usr/local/bin/gunicorn wsgi --bind localhost:2000
directory=/root/zxd-bbs
autostart=true
autorestart=true
