# zxd-bbs.conf
[program:web20]
command=/usr/local/bin/gunicorn wsgi --bind localhost:2000
directory=/var/www/web23
autostart=true
autorestart=true
